# systemd/modbus_tcp_proxy.service
[Unit]
Description=Modbus TCP Proxy for GoodWe inverter
After=network-online.target
Wants=network-online.target

[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/rpi-admin-ui/modbus_tcp_proxy.py
WorkingDirectory=/opt/rpi-admin-ui
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
