# systemd/rpi-admin-ui.service
[Unit]
Description=RPi Admin UI
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/rpi-admin-ui
ExecStart=/usr/bin/python3 /opt/rpi-admin-ui/app.py
Restart=on-failure
RestartSec=5
User=pi
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
